# daemon/systemd/aliendalvik-control.service
[Unit]
Description=Aliendalvik Control
After=dbus.socket
Requires=dbus.socket

[Service]
Type=dbus
ExecStart=/usr/bin/aliendalvik-control
EnvironmentFile=-/var/lib/environment/aliendalvik-control/*.conf
BusName=org.coderus.aliendalvikcontrol
Restart=on-failure
RestartSec=15
